# openclaw-trading.service
[Unit]
Description=OpenClaw Trading System
After=network.target

[Service]
Type=simple
User=cool
WorkingDirectory=/home/cool/ai-trading-system
ExecStart=/home/cool/ai-trading-system/.venv/bin/python -m src.main
Restart=always
RestartSec=5
Environment=PYTHONPATH=/home/cool/ai-trading-system
Environment=TRADING_ENV=production
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target