# deploy/systemd/openclaw-live-stability-monitor.timer
[Unit]
Description=Run OpenClaw live stability monitor every 2 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=2h
AccuracySec=1min
Persistent=true
Unit=openclaw-live-stability-monitor.service

[Install]
WantedBy=timers.target
