# deploy/systemd/openclaw-health-audit.service
[Unit]
Description=OpenClaw Full System Audit
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/cool/ai-trading-system
EnvironmentFile=-/home/cool/ai-trading-system/.env
ExecStart=/home/cool/ai-trading-system/scripts/run_scheduled_health_suite.sh audit
SuccessExitStatus=2

[Install]
WantedBy=multi-user.target
