# scripts/openclaw-trading.service
[Unit]
Description=OpenClaw Trading System
After=network.target

[Service]
Type=simple
User=cool
Group=cool
WorkingDirectory=/home/cool/.openclaw-trading
Environment="PATH=/home/cool/.openclaw-trading:/usr/local/bin:/usr/bin"
ExecStart=/usr/bin/python3 -m src.main
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
