# scripts/openclaw-trading.service
[Unit]
Description=OpenClaw Trading System
After=network.target

[Service]
Type=simple
User=cool
Group=cool
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/home/cool/.openclaw-trading/scripts/start-openclaw-trading.sh
ExecStop=/home/cool/.openclaw-trading/scripts/stop-openclaw-trading.sh
WorkingDirectory=/home/cool/.openclaw-trading
Environment="PATH=/home/cool/.openclaw-trading:/usr/local/bin:/usr/bin"

[Timer]
OnCalendar=*-*-* **:00:00
Unit=openclaw-trading.service

[Install]
WantedBy=multi-user.target
